# config/systemd/user/bt-toggle.service
[Unit]
Description=Toggle Bluetooth and connect Q30

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bt-toggle
