# config/systemd/user/bw-serve.service
[Unit]
Description=Bitwarden CLI REST API
After=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'cat ${XDG_DATA_HOME:-$HOME/.local/share}/bw/session'
ExecStart=/bin/bash -c 'export BW_SESSION=$(cat ${XDG_DATA_HOME:-$HOME/.local/share}/bw/session) && exec bw serve --port 8087 --hostname 127.0.0.1'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
