# config/systemd/user/screen-off-toggle.service
[Unit]
Description=Toggle screen off and mute notifications

[Service]
Type=oneshot
ExecStart=%h/.local/bin/screen-off-toggle
