# config/systemd/user/screen-off-watcher.service
[Unit]
Description=Watch for screen wake and restore notifications

[Service]
Type=simple
ExecStart=%h/.local/bin/screen-off-watcher
Restart=no
