# config/systemd/usb-hub-bt-on.service
[Unit]
Description=Turn Bluetooth on when USB hub reconnects

[Service]
Type=oneshot
ExecStart=/home/nicholas/.local/bin/usb-hub-bt-on
