# config/systemd/usb-hub-bt-off.service
[Unit]
Description=Turn Bluetooth off on USB hub disconnect

[Service]
Type=oneshot
ExecStart=/home/nicholas/.local/bin/usb-hub-bt-off
