# config/systemd/ollama.service.d/override.conf
[Service]
Environment="OLLAMA_MODELS=/mnt/data/ollama/"
